# src/resources/info.service
[Unit]
Description=AKTIN Update-Info Service
After=network.target

[Service]
ExecStart=apt update

[Install]
WantedBy=default.target

